# udev/usb-wifi/50-usb-wifi.rules
KERNEL=="<IF-NAME>", ATTR{address}=="<MAC-ADDR>", RUN+="/opt/udev-usb-wifi.sh %k"
